# ha_governance/rootfs/etc/s6-overlay/s6-rc.d/ha_governance/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: HA Governance
# Handles service shutdown
# ==============================================================================

bashio::log.info "HA Governance service stopped"
