# ha_governance/rootfs/etc/s6-overlay/s6-rc.d/ha_governance/run
#!/usr/bin/with-contenv bash
set -e

echo "[HA Governance] Starting service..."

cd /app
exec python3 -u main.py
